# dockerhosts.service
[Unit]
Description=docker-hosts

[Service]
ExecStart=/usr/bin/dockerhosts

[Install]
WantedBy=multi-user.target
